# Systemd/tunsnx.service
[Unit]
Description=Log test with timer

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/tunsnx.sh

[Install]
WantedBy=multi-user.target
